# firefox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -ac
priority=10
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.log

[program:openbox]
command=/usr/bin/openbox
environment=DISPLAY=":1",HOME="/root"
priority=20
autorestart=true

[program:install-foxmcp]
command=/opt/install-foxmcp.sh
priority=25
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/var/log/install-foxmcp.log
stderr_logfile=/var/log/install-foxmcp.log

[program:firefox]
command=/usr/bin/firefox-esr --no-remote --profile /root/.mozilla/firefox/default-esr
environment=DISPLAY=":1",HOME="/root"
priority=30
autorestart=true
startsecs=5
stdout_logfile=/var/log/firefox.log
stderr_logfile=/var/log/firefox.log

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -rfbport 5900 -noxdamage
priority=40
autorestart=true
stdout_logfile=/var/log/x11vnc.log

[program:novnc]
command=/opt/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:6080
priority=50
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.log
